# PGS/config/cfengine/inputs/global/cfgfiles/profile.d/holoSeis.csh
if ( ! $?HOLOSEIS_ROOT_DIR ) then
   setenv HOLOSEIS_ROOT_DIR /tps/holoSeis/Linux64/
   set path = ( $path $HOLOSEIS_ROOT_DIR )
endif

